# supervisord.d/hinemos-web.ini
[program:hinemos-web]
user = root
command=/opt/hinemos_web/bin/tomcat_start.sh -Wq


