# supervisord.d/hinemos-manager.ini
[program:hinemos-manager]
user = root
command=/opt/hinemos/bin/jvm_start.sh -W



